# local/sv/rivertile/run
#!/bin/sh
exec rivertile -view-padding 4 -outer-padding 4
