# local/sv/swayidle/run
#!/bin/sh
exec swayidle -w \
	timeout 300 "waylock.sh" \
	timeout 600 "wlopm --off '*'" resume "wlopm --on '*'"
